# cron/crontab.cron
###################################################################################
# safety nets (checking that processes are up)
###################################################################################

*/3 * * * *  sleep 1;   source ~/.crypto_bot_profile;$CRYPTO_BOT_BIN_SCRIPT_DIR/kraken_bridge_orderbook_btc_usdt_1_safety_net.sh    # every 3 minutes
*/3 * * * *  sleep 3;   source ~/.crypto_bot_profile;$CRYPTO_BOT_BIN_SCRIPT_DIR/kraken_bridge_orderbook_btc_usdt_2_safety_net.sh    # every 3 minutes
*/3 * * * *  sleep 6;   source ~/.crypto_bot_profile;$CRYPTO_BOT_BIN_SCRIPT_DIR/kraken_bridge_orderbook_btc_usdt_3_safety_net.sh    # every 3 minutes
*/3 * * * *  sleep 9;   source ~/.crypto_bot_profile;$CRYPTO_BOT_BIN_SCRIPT_DIR/kraken_bridge_orderbook_btc_usdt_4_safety_net.sh    # every 3 minutes
*/3 * * * *  sleep 12;  source ~/.crypto_bot_profile;$CRYPTO_BOT_BIN_SCRIPT_DIR/kraken_bridge_orderbook_btc_usdt_5_safety_net.sh    # every 3 minutes
*/3 * * * *  sleep 15;  source ~/.crypto_bot_profile;$CRYPTO_BOT_BIN_SCRIPT_DIR/kraken_bridge_orderbook_btc_usdt_6_safety_net.sh    # every 3 minutes
*/3 * * * *  sleep 18;  source ~/.crypto_bot_profile;$CRYPTO_BOT_BIN_SCRIPT_DIR/kraken_bridge_orderbook_btc_usdt_7_safety_net.sh    # every 3 minutes
*/3 * * * *  sleep 21;  source ~/.crypto_bot_profile;$CRYPTO_BOT_BIN_SCRIPT_DIR/kraken_bridge_orderbook_btc_usdt_8_safety_net.sh    # every 3 minutes
*/3 * * * *  sleep 24;  source ~/.crypto_bot_profile;$CRYPTO_BOT_BIN_SCRIPT_DIR/kraken_bridge_orderbook_btc_usdt_9_safety_net.sh    # every 3 minutes
*/3 * * * *  sleep 27;  source ~/.crypto_bot_profile;$CRYPTO_BOT_BIN_SCRIPT_DIR/kraken_bridge_orderbook_btc_usdt_10_safety_net.sh   # every 3 minutes

*/4 * * * *  sleep 22;  source ~/.crypto_bot_profile;$CRYPTO_BOT_BIN_SCRIPT_DIR/kraken_persister_orderbook_safety_net.sh            # every 4 minutes

###################################################################################
# monitoring
###################################################################################

*/1 * * * *  sleep 7;   source ~/.crypto_bot_profile;$CRYPTO_BOT_BIN_SCRIPT_DIR/jms_queue_monitor.sh            # every minute
*/5 * * * *  sleep 41;  source ~/.crypto_bot_profile;$CRYPTO_BOT_BIN_SCRIPT_DIR/data_age_monitor.sh             # every 5 minutes

###################################################################################
# cleanup
###################################################################################

*/2 * * * *  sleep 51;  source ~/.crypto_bot_profile;$CRYPTO_BOT_BIN_SCRIPT_DIR/data_cleaner.sh                 # every 2 minutes
*/32 * * * *  sleep 57; source ~/.crypto_bot_profile;find $CRYPTO_BOT_LOG_DIR/* -mtime +7 -type f -delete       # prune logs older than 7 days - every 32 mins

###################################################################################
# restart all daemons
###################################################################################

# 14 */5 * * *  sleep 12; source ~/.crypto_bot_profile;$CRYPTO_BOT_BIN_SCRIPT_DIR/restart_processes.sh   # at minute 14 past every 5th hour
